# docker/python/crontab
* 8 * * *         root  sh /usr/local/bin/task.sh

