# docker/python/crontab
* 8 * * *         root  sh /usr/local/bin/task.sh
*/2 * * * *         root  sh /usr/local/bin/hello.sh

