# docker/python/example.crontab
5 7 * * * /usr/local/bin/task01.sh
30 9 * * * /usr/local/bin/task02.sh
50 11 * * * /usr/local/bin/task03.sh
